# build/conf/crontab
*/1 * * * * apache php /usr/local/cacti/poller.php --force > /dev/null 2>&1
